# other_stuff/profile.d/grails.sh
export GRAILS_HOME=/opt/grails/latest
export PATH=${GRAILS_HOME}/bin:${PATH}
